# playbook-awx-artifactory-download-and-deploy-main/roles/download_and_deploy_application/templates/app.service.j2
[Unit]
Description={{ application_name }}
After=syslog.target

[Service]
User=root
ExecStart=/bin/java -jar {{ application_directory }}/{{ application_name }}-{{ application_version }}.{{ application_extension }}
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
